# xoa-credentials.service
[Unit]
Description=XOA First Boot - Admin credentials provisioning (phase 2)
# Runs AFTER xo-server is supposed to be up
After=network-online.target xo-server.service
Wants=network-online.target
ConditionPathExists=!/var/lib/xoa-credentials.done

[Service]
Type=oneshot
ExecStart=/opt/xoa-credentials.sh
RemainAfterExit=yes
# Retry a few times in case xo-server is slow
Restart=on-failure
RestartSec=15
StartLimitIntervalSec=300
StartLimitBurst=5

[Install]
WantedBy=multi-user.target
