# systemd/xoa-credentials.service
[Unit]
Description=XOA First Boot - Admin credentials provisioning (phase 2)
After=network-online.target xo-server.service
Wants=network-online.target
ConditionPathExists=!/var/lib/xoa-credentials.done

[Service]
Type=oneshot
ExecStart=/root/xoa-credentials.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
